# icecast-streamer.service
[Unit]
Description=Icecast Streaming Encoder
After=network.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/icecast-streamer
ExecStart=/usr/bin/python3 /home/pi/icecast-streamer/streamer.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

# Security settings
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
